# com.citadel.voice.plist

<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.citadel.voice</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/homelab/voice/.venv/bin/uvicorn</string>
        <string>main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8001</string>
        <string>--log-level</string>
        <string>info</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/homelab/voice</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/homelab</string>
        <key>WHISPER_BACKEND_URL</key>
        <string>http://127.0.0.1:8003</string>
        <key>HEAD_TRIM_MS</key>
        <string>50</string>
        <key>BACKEND_TIMEOUT</key>
        <string>30.0</string>
        <!-- VOICE_BEARER_TOKEN: add here if auth is required -->
        <!-- <key>VOICE_BEARER_TOKEN</key> -->
        <!-- <string>your-token-here</string> -->
    </dict>

    <key>UserName</key>
    <string>homelab</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/homelab/voice/logs/voice.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/homelab/voice/logs/voice.err</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>

